# dot_config/systemd/user/vaultwarden-backup.service
[Unit]
Description=Vaultwarden Backup

[Service]
Type=oneshot
ExecStart=/bin/bash %h/Scripts/backup-vaultwarden.sh
